# dotfiles/rslsync/.config/autostart/resilio-sync.desktop
[Desktop Entry]
Type=Application
Name=Resilio Sync
Exec=rslsync
Icon=resilio-sync
Comment=Start Resilio Sync on login
