# dotfiles/syncthing/.config/autostart/syncthing.desktop
[Desktop Entry]
Type=Application
Name=Syncthing
Exec=syncthing
Icon=syncthing
Comment=Start Syncthing on login
