# dotfiles/syncthing/.config/autostart/syncthing.desktop
[Desktop Entry]
Name=Syncthing
Type=Application
Exec=syncthing serve --no-browser --logfile=default
Icon=syncthing
Comment=Start syncthing on login
